# examples/random-splash.service
[Unit]
Description=Get a new wallpaper from Unsplash

[Service]
Type=oneshot
ExecStart=/home/example/bin/getrandomwallpaper.sh
